# com.ali.aave-autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ali.aave-autopush</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ali/Documents/python/aave_v3.1/auto_push.sh</string>
    </array>

    <!-- Run every day at 23:00 (11:00 PM) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ali/Documents/python/aave_v3.1/.autopush.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ali/Documents/python/aave_v3.1/.autopush.err.log</string>
</dict>
</plist>
